# crontab
0 9 * * 0-4 cd /app && /bin/bash -c "/usr/local/bin/python3 scheduler.py" >> /var/log/cron.log 2>&1 
